# backend/systemd/autoscript-web.service
[Unit]
Description=Autoscript Web (TanStack Start panel server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/autoscript/agent.env
Environment=HOST=127.0.0.1
Environment=NODE_ENV=production
WorkingDirectory=/opt/autoscript
ExecStart=/bin/sh -c 'PORT=${WEB_INTERNAL_PORT} exec /usr/bin/node /opt/autoscript/backend/web-server.mjs'
Restart=on-failure
RestartSec=3
User=root
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
