# backend/systemd/autoscript-agent.service
[Unit]
Description=Autoscript Panel Agent (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/autoscript/agent.env
WorkingDirectory=/opt/autoscript/backend/agent
ExecStart=/opt/autoscript/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8088
Restart=on-failure
RestartSec=3
# Agent shells out to scripts that need root (tc, useradd, systemctl).
User=root

[Install]
WantedBy=multi-user.target
